# deploy/supervisor/umanie.conf
; Supervisor config for Umanie Homes ERP
; Place this file at: /etc/supervisor/conf.d/umanie.conf
; Then run: sudo supervisorctl reread && sudo supervisorctl update

[program:umanie-gunicorn]
command=/home/ubuntu/umanie/venv/bin/gunicorn config.asgi:application
    -k uvicorn.workers.UvicornWorker
    --workers 4
    --bind 127.0.0.1:8000
    --timeout 120
    --log-level info
directory=/home/ubuntu/umanie/backend
user=ubuntu
environment=DJANGO_SETTINGS_MODULE="config.settings.production"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/umanie/gunicorn.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:umanie-celery]
command=/home/ubuntu/umanie/venv/bin/celery -A config worker -l info -c 4
directory=/home/ubuntu/umanie/backend
user=ubuntu
environment=DJANGO_SETTINGS_MODULE="config.settings.production"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/umanie/celery.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:umanie-celery-beat]
command=/home/ubuntu/umanie/venv/bin/celery -A config beat -l info
directory=/home/ubuntu/umanie/backend
user=ubuntu
environment=DJANGO_SETTINGS_MODULE="config.settings.production"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/umanie/celery-beat.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2

[group:umanie]
programs=umanie-gunicorn,umanie-celery,umanie-celery-beat
